# .config/systemd/user/dunst-switcher.service
[Unit]
Description=Dunst Switcher

[Service]
ExecStart=/bin/bash ${HOME}/.local/bin/dunst-switcher

[Install]
WantedBy=default.target
