# .root/etc/systemd/system/openvpn-client@.service.d/netns.conf
[Unit]
Requires=netns-vpn.service
After=netns-vpn.service

[Service]
AmbientCapabilities=CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN
