# .root/etc/NetworkManager/dispatcher.d/rm-lease
#!/bin/bash
if [[ $NM_DISPATCHER_ACTION = down ]]; then
    rm -f "/var/lib/NetworkManager/internal-$CONNECTION_UUID-$DEVICE_IP_IFACE.lease"
fi
